# configs/autostart/.config/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Name=loginctl
Exec=/home/jan/.loginctl
Comment=Run loginctl script on login
X-GNOME-Autostart-enabled=true
